# cookbook-api/cookbook-api.service
[Unit]
Description=Cookbook API Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/cookbook-api
ExecStart=/opt/cookbook-api/start-cookbook.sh
Restart=on-failure
RestartSec=5s
Environment=ASPNETCORE_ENVIRONMENT=Production

# Adjusted resource limits for t3.small
CPUQuota=90%
MemoryHigh=1.5G
MemoryMax=1.8G
TasksMax=200
LimitNPROC=200
LimitNOFILE=8192

# Run cleanup script as root on stop
ExecStopPost=/bin/sh -c 'sudo /opt/cookbook-api/cleanup-playwright.sh'

[Install]
WantedBy=multi-user.target
